# deploy/smartgrader.service
[Unit]
Description=SmartGrader Exam Platform
After=network.target

[Service]
User=smartgrader
Group=smartgrader
WorkingDirectory=/opt/smartgrader
EnvironmentFile=/opt/smartgrader/.env
ExecStart=/opt/smartgrader/venv/bin/gunicorn -c deploy/gunicorn.conf.py "app:create_app('production')"
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
